# infra/systemd/nexus_trader.service
[Unit]
Description=Nexus Trader V25 Production Stack
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/nexus-trader
EnvironmentFile=/opt/nexus-trader/infra/production.env
ExecStart=/usr/bin/python3 /opt/nexus-trader/scripts/start_production.py --service all --loop
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
